# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true

[program:deseq2shiny]
command=/opt/conda/envs/v_deseq2/bin/R -e "shiny::runApp('/srv/shiny-server/deseq2shiny/src', port=4001, host='0.0.0.0')"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:ATACseqQCShniy]
command=/opt/conda/envs/v_atacqc/bin/R -e "shiny::runApp('/srv/shiny-server/ATACseqQCShniy/src', port=4002, host='0.0.0.0')"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:ClusterProfShinyGSEA]
command=/opt/conda/envs/v_enrichGSEA/bin/R -e "shiny::runApp('/srv/shiny-server/ClusterProfShinyGSEA/src', port=4003, host='0.0.0.0')"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:ClusterProfShinyORA]
command=/opt/conda/envs/v_enrichORA/bin/R -e "shiny::runApp('/srv/shiny-server/ClusterProfShinyORA/src', port=4004, host='0.0.0.0')"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:GeneCountMerger]
command=/opt/conda/envs/v_countmerger/bin/R -e "shiny::runApp('/srv/shiny-server/GeneCountMerger/src', port=4005, host='0.0.0.0')"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:SeuratV5Shiny]
command=/opt/conda/envs/v_seuratv5/bin/R -e "shiny::runApp('/srv/shiny-server/SeuratV5Shiny/src', port=4006, host='0.0.0.0')"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:animalcules]
command=/opt/conda/envs/v_animalcules/bin/R -e "shiny::runApp('/srv/shiny-server/animalcules/src', port=4007, host='0.0.0.0')"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:dada2Shiny]
command=/opt/conda/envs/v_dada2/bin/R -e "shiny::runApp('/srv/shiny-server/dada2Shiny/src', port=4008, host='0.0.0.0')"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:monocle3]
command=/opt/conda/envs/v_monocle3/bin/R -e "shiny::runApp('/srv/shiny-server/monocle3/src', port=4009, host='0.0.0.0')"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:DEBrowser]
command=/opt/conda/envs/v_debrowser/bin/R -e "shiny::runApp('/srv/shiny-server/DEBrowser/src/R', port=4010, host='0.0.0.0')"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:mergeFPKMs]
command=/opt/conda/envs/v_mergefpkms/bin/R -e "shiny::runApp('/srv/shiny-server/mergeFPKMs/src', port=4011, host='0.0.0.0')"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:tsar_nasqar]
command=/opt/conda/envs/v_tsar/bin/R -e "shiny::runApp('/srv/shiny-server/tsar_nasqar/src', port=4012, host='0.0.0.0')"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0













